# crontab
0 8 * * * ~/code/jobs-notify/jobs-notify.sh