# crontab
30 08 10 06 * ~/code/jobs-notify/jobs-notify.sh